# gemini-api.service
[Unit]
Description=Gemini Browser OpenAI-compatible API server
# Chrome runs on a virtual X display (Xvfb) via xvfb-run — no real display needed,
# so the service is fully background and independent of any login session.

[Service]
Type=simple
WorkingDirectory=/home/b/google_api
# Clear a stale profile lock left by a previous hard crash (ignore if absent).
ExecStartPre=-/bin/rm -f /home/b/google_api/gemini_profile/SingletonLock /home/b/google_api/gemini_profile/SingletonSocket /home/b/google_api/gemini_profile/SingletonCookie
# xvfb-run gives Chrome an invisible virtual display (auto-picked server number).
ExecStart=/usr/bin/xvfb-run -a /home/b/.pyenv/versions/3.12.1/bin/python3 /home/b/google_api/server.py
# nodriver detaches Chrome from the cgroup; reap any survivors on stop/restart.
ExecStopPost=-/usr/bin/pkill -TERM -f "user-data-dir=./gemini_profile"
Restart=always
RestartSec=5
KillMode=control-group
TimeoutStopSec=20

[Install]
WantedBy=default.target
