# gemini-api.service
[Unit]
Description=Gemini Browser OpenAI-compatible API server
# Chrome runs on a virtual X display (Xvfb) via xvfb-run — no real display needed,
# so the service is fully background and independent of any login session.

[Service]
Type=simple
WorkingDirectory=/home/eben/Downloads/google_api
# Where generated images are saved (change this to set your storage path).
Environment=GEMINI_IMAGE_DIR=/home/eben/Pictures/gemini
# Base URL the saved images are served under (used in returned links).
Environment=GEMINI_PUBLIC_URL=http://localhost:8081
# Clear stale profile locks left by a previous hard crash (ignore if absent).
ExecStartPre=-/bin/rm -f /home/eben/Downloads/google_api/gemini_profile/SingletonLock /home/eben/Downloads/google_api/gemini_profile/SingletonSocket /home/eben/Downloads/google_api/gemini_profile/SingletonCookie
ExecStartPre=-/bin/rm -f /home/eben/Downloads/google_api/chatgpt_profile/SingletonLock /home/eben/Downloads/google_api/chatgpt_profile/SingletonSocket /home/eben/Downloads/google_api/chatgpt_profile/SingletonCookie
# xvfb-run gives Chrome an invisible virtual display (auto-picked server number).
ExecStart=/usr/bin/xvfb-run -a /home/eben/Downloads/google_api/venv/bin/python /home/eben/Downloads/google_api/server.py
# nodriver detaches Chrome from the cgroup; reap any survivors on stop/restart
# (both provider profiles: gemini_profile / chatgpt_profile).
ExecStopPost=-/usr/bin/pkill -TERM -f "user-data-dir=./gemini_profile"
ExecStopPost=-/usr/bin/pkill -TERM -f "user-data-dir=./chatgpt_profile"
Restart=always
RestartSec=5
KillMode=control-group
TimeoutStopSec=20

[Install]
WantedBy=default.target
